# launchd/com.construction-ai.git-autosave.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.construction-ai.git-autosave</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/rodrigo/construction-ai/scripts/git_autosave.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/construction-ai-git-autosave.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/construction-ai-git-autosave.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
